# container/home/.supervisor_conf.d/xvfb.conf
[program:xvfb]
command=bash -c 'export DISPLAY=:0 && Xvfb :0 -screen 0 1024x768x24'
autostart=true
autorestart=false
stderr_logfile=/var/log/xvfb.err.log
priority=997
